# roles/vmsingle/templates/victoriametrics.service.j2
# Ansible managed: Do NOT edit this file manually!

[Unit]
Description=VictoriaMetrics service
After=network.target

[Service]
Type=simple
LimitNOFILE={{ vmsingle_max_open_files }}
User={{ vmsingle_system_user }}
Group={{ vmsingle_system_group }}
ExecStart=/usr/local/bin/victoria-metrics-prod {% if (vmsingle_service_envflag_enabled | bool) %} -envflag.enable {% endif %} {% for flag, flag_value in vmsingle_service_args.items() %}-{{ flag }}={{ flag_value }} {% endfor %}

SyslogIdentifier=victoriametrics
Restart=always

PrivateTmp=yes
ProtectHome=yes
NoNewPrivileges=yes

ProtectSystem=full

{% if vmsingle_systemd_version | int >= 232 %}
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=yes
{% endif %}

{% if vmsingle_service_envflag_enabled | bool %}
{% for v in vmsingle_service_envflag_data %}
Environment="{{ v }}"
{% endfor %}

EnvironmentFile={{ vmsingle_service_envflag_file }}
{% endif %}

[Install]
WantedBy=multi-user.target
